# etc/vc3-master.service
[Unit]
Description=VC3 Master

[Service]
ExecStart=/usr/bin/vc3-master -c /etc/vc3/master.conf
Restart=on-failure

[Install]
WantedBy=multi-user.target

